# twitter-globe.conf
description 'twitter-globe'
author 'psi <psi@ledyba.org>'
start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

setuid psi
setgid psi

respawn

exec /opt/run/twitter-globe/twitter-globe
