# twitter-globe.conf
[program:twitter-globe]
command=/opt/run/twitter-globe/twitter-globe
directory=/opt/run/twitter-globe
user=psi
