# twitter-globe.conf
[program:twitter-globe]
command=/opt/run/twitter-globe/twitter-globe -port 8083
directory=/opt/run/twitter-globe
user=psi
